# lib/setuputil/files/god.conf
description "god - by bitc-configurator"

start on (local-filesystems and net-device-up IFACE!=lo)
# start on runlevel [2]
stop on runlevel [016]

console log

exec /usr/local/bin/god -c /etc/bitc/god/bitc.god
expect fork

respawn
